# dots/.config/systemd/user/hyprland-session.service
[Unit]
Description=Hyprland graphical session
BindsTo=graphical-session.target
Before=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/bin/sleep infinity
RemainAfterExit=no
